# repo-refresh/repo-refresh.service
[Unit]
Description=Check for package update
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/usr/libexec/repo-refresh.sh
